# isv:perlur:CardanoOS:MicroOS/systemd-mini/.osc/80-acpi-container-hotplug.rules
# ACPI0004 container offline for Huawei Kunlun
# do not edit this file, it will be overwritten on update

SUBSYSTEM=="container", ACTION=="change", ENV{EVENT}=="offline", ATTR{online}=="1", ATTR{[dmi/id]sys_vendor}=="Huawei", ATTR{[dmi/id]product_name}=="9008|9016|9032", DEVPATH=="*/ACPI0004:??", \
RUN+="/usr/bin/sh -c ' \
        /usr/bin/find -L /sys/$$env{DEVPATH}/firmware_node/*/physical_node* -maxdepth 1 -name online | \
        while read line; do \
                if [ $$(/usr/bin/cat $$line) -eq 1 ]; then \
                        /usr/bin/echo 0 > $$line; \
                fi \
        done; \
        /usr/bin/echo 0 > /sys/$$env{DEVPATH}/online; '"
